# launchd/com.charles.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.charles.morning-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/home/charles/.venv/bin/python</string>
        <string>/Users/home/charles/scripts/morning_brief.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/home/charles</string>

    <!-- Fire once daily at 06:00 local time (host TZ = America/New_York,
         labeled EST per John's preference). LaunchAgents use local time. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/home/charles/logs/morning_brief.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/home/charles/logs/morning_brief.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/home/charles/.venv/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
