# launchd/com.charles.caffeinate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.charles.caffeinate</string>

    <!-- caffeinate -i prevents idle sleep; -s prevents system sleep when on AC.
         Display can still sleep. KeepAlive ensures it respawns if killed. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-s</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/home/charles/logs/caffeinate.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/home/charles/logs/caffeinate.err.log</string>
</dict>
</plist>
